# etc/init.d/dataset-repo
#!/sbin/openrc-run

name="Relational Dataset Repository"

: ${command_user:="web:web"}
: ${directory:="/opt/dataset-repo"}
: ${error_log="/var/log/dataset-repo.log"}
: ${output_log="/var/log/dataset-repo.log"}

command="/usr/bin/node"
command_args="$directory/src/server"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

export NODE_ENV=production

depend() {
	use net
	need mariadb
}

start_pre() {
	if [ "$error_log" ]; then
		checkpath -f -m 640 -o "$command_user" "$error_log" || return 1
	fi
	if [ "$output_log" ]; then
		checkpath -f -m 640 -o "$command_user" "$output_log" || return 1
	fi
	return 0
}
